# crontab
*/1 * * * * ~/home/r4d3k/Pulpit/Zadania FrontStep/Level_5/NewLocalRep/pushScript.sh
